# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:telegram-bot-api]
command=telegram-bot-api
    --api-id=%(ENV_TELEGRAM_API_ID)s
    --api-hash=%(ENV_TELEGRAM_API_HASH)s
    --local
    --dir=/var/lib/tgapi
    --temp-dir=/tmp/tgapi
    --http-port=8081
priority=1
autorestart=true
stdout_logfile=/var/log/supervisor/tgapi.log
stderr_logfile=/var/log/supervisor/tgapi.log
startsecs=5

[program:bot]
command=python -u /app/bot.py
priority=2
autorestart=true
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot.log
startsecs=10
